# systemd_files/pzkids.timer
[Unit]
Description=PZKids Timer

[Timer]
OnCalendar=Mon..Fri 16:10
Persistent=true

[Install]
WantedBy=default.target
